# fenix/archives/filesystem/common/etc/update-motd.d/10-help-text
#!/bin/sh

printf "\n"
printf " * Website:        https://jevois.org\n"
printf " * Documentation:  https://jevois.org/doc\n"
printf " * Forum:          https://jevois.org/qa\n"
printf "\n"
